# rootfs/etc/profile.d/global.sh
#!/bin/bash

# Default editor
export EDITOR='/bin/nano'

# Default aliases
alias l='ls -lA --color=auto'
alias ll='ls -l --color=auto'
alias ls='ls --color=auto'
alias tree='tree --dirsfirst -C'
